# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:nginx]
command=nginx -c /app/nginx.conf -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=10

[program:personaplex]
command=/app/moshi/.venv/bin/python -m moshi.server --host 0.0.0.0 --port 8999
directory=/app/moshi
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=20
startsecs=30
startretries=3

[program:qwen3tts]
command=/app/qwen3-tts/.venv/bin/python -m api.main
directory=/app/qwen3-tts
environment=HOST="0.0.0.0",PORT="8880",WORKERS="1",PYTHONPATH="/app/qwen3-tts",TTS_BACKEND="official",TTS_MODEL_NAME="Qwen/Qwen3-TTS-12Hz-1.7B-Base",TTS_WARMUP_ON_START="false",NVIDIA_VISIBLE_DEVICES="all"
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=20
startsecs=30
startretries=3
